# bdbot/misc/runbdbot.service
[Unit]
Description=BDBot

[Service]
WorkingDirectory={PACKAGE_DIRECTORY}
User={USER}
ExecStart=/home/{USER}/.local/bin/runbdbot
Type=forking
Restart=always

[Install]
WantedBy=multi-user.target
